# roles/experiment-job-start/templates/client.service.j2
[Unit]
Description=Client Service

# TODO define how to start the client service and with which arguments

[Service]
Type=simple
WorkingDirectory={{ exp_working_dir }}
ExecStart={{ exp_code_dir }}/demo/client.py -a1 {{ exp_config.client.arg1 }} -a2 {{ exp_config.client.arg2 }}

StandardOutput=file:{{ exp_working_dir }}/results/server_stdout.log
StandardError=file:{{ exp_working_dir }}/results/server_stderr.log

[Install]
WantedBy=multi-user.target