# roles/experiment-job-start/templates/client.service.j2
[Unit]
Description=Client Service

# TODO: define how to start the client service and with which arguments
#       the run config is available in the variable `exp_run_config` to pass as parameters
#       or in the working directory there is a file called `config.json` with the run config

# Note:
# - 'host_ips' is a variable with a dictionary containing the private IPs of
#    other hosts belonging to this experiment.
#    host_ips = { host_type: [ 10.100.0.xy, 10.100.0.yz, ... ], ... }

[Service]
Type=simple
WorkingDirectory={{ exp_working_dir }}
# TODO: add something like this if it needs to wait for another machine
# ExecStartPre=/bin/bash -c 'until host example.com; do sleep 1; done'
ExecStart=/bin/sh -c '{{ exp_code_dir }}/.venv/bin/python {{ exp_code_dir }}/demo/client.py -config config.json -a1 {{ exp_run_config.client.arg1 }} -a2 {{ exp_run_config.client.arg2 }}'
User=ubuntu

StandardOutput=file:{{ exp_working_dir }}/results/client_stdout.log
StandardError=file:{{ exp_working_dir }}/results/client_stderr.log

[Install]
WantedBy=multi-user.target
