# roles/experiment-job-start/templates/server.service.j2
[Unit]
Description=Server Service

# TODO define how to start the server service and with which arguments

[Service]
Type=simple
WorkingDirectory={{ exp_working_dir }}
ExecStart={{ exp_code_dir }}/demo/server.py -a3 {{ exp_config.server.arg3 }} -a4 {{ exp_config.server.arg4 }}

StandardOutput=file:{{ exp_working_dir }}/results/server_stdout.log
StandardError=file:{{ exp_working_dir }}/results/server_stderr.log

[Install]
WantedBy=multi-user.target