# roles/experiment-job-start/templates/single.service.j2
[Unit]
Description=Single Service

# TODO: define how to start the service and with which arguments
#       the run config is available in the variable `exp_run_config` to pass as parameters
#       or in the working directory there is a file called `config.json` with the run config

# Note:
# - 'host_ips' is a variable with a dictionary containing the private IPs of
#    other hosts belonging to this experiment.
#    host_ips = { host_type: [ 10.100.0.xy, 10.100.0.yz, ... ], ... }

[Service]
Type=simple
WorkingDirectory={{ exp_working_dir }}
ExecStart=/usr/bin/echo "{ \"rt_mean\": $((1 + $RANDOM % 100)).$(($RANDOM % 10)), \"rt_std\": $((1 + $RANDOM % 10)).$(($RANDOM % 10)) }"
User=ubuntu

StandardOutput=file:{{ exp_working_dir }}/results/client.json
StandardError=file:{{ exp_working_dir }}/results/client_stderr.log

[Install]
WantedBy=multi-user.target
